# scripts/com.daedalus.server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.daedalus.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/uv</string>
    <string>run</string>
    <string>daedalus</string>
    <string>serve</string>
    <string>--model</string><string>mlx-community/Qwen3.5-9B-MLX-4bit</string>
    <!-- Keep the secret OUT of this plist (launchd plists are world-readable
         and end up in backups). Store it in a chmod-600 file instead:
           echo "your-secret" > ~/.config/daedalus/api-key && chmod 600 ~/.config/daedalus/api-key -->
    <string>--api-key-file</string><string>/Users/YOU/.config/daedalus/api-key</string>
    <string>--max-active-memory-gb</string><string>12</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/daedalus.log</string>
  <key>StandardErrorPath</key><string>/tmp/daedalus-error.log</string>
</dict>
</plist>
